# config/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D -o PasswordAuthentication=no
autorestart=true
